# vlog.service
[Unit]
Description=VLog Autonomous Life Logger Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/kafka/vlog
ExecStart=/usr/bin/cargo run --release -- monitor
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
