# vlog.service
[Unit]
Description=VRChat Auto-Diary Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/kafka/projects/vlog
ExecStart=/home/kafka/.local/bin/uv run python -m src.main
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
