# systemd/vlog-daily.timer
[Unit]
Description=Run VLog Daily Processing at 5AM JST

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
Unit=vlog-daily.service

[Install]
WantedBy=timers.target
