# systemd/vlog-daily.service
[Unit]
Description=VLog Daily Processing Service

[Service]
Type=oneshot
User=kafka
WorkingDirectory=/home/kafka/vlog
ExecStart=/home/kafka/vlog/target/release/vlog-rs evaluate --date yesterday
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=vlog-daily
